# claude-code-shared/launchd/com.ericlingren.cc-usage-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ericlingren.cc-usage-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/eric/.dotfiles/claude-code-shared/scripts/weekly-usage-report.sh</string>
  </array>
  <!-- Mondays at 09:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/eric/.cache/cc-usage-reports/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/eric/.cache/cc-usage-reports/launchd.err.log</string>
</dict>
</plist>
